# etc/s6-overlay/s6-rc.d/fix-user-ids/up
#!/command/with-contenv bash

/var/fixes/set_user.sh
