# etc/s6-overlay/s6-rc.d/start-ssh-server/run
#!/command/with-contenv bash

exec 2>&1
exec s6-setuidgid $USER /usr/sbin/sshd -f /ssh/sshd_config/config -D -e
