# etc/s6-overlay/s6-rc.d/start-ssh-server/run
#!/command/with-contenv bash

echo "Starting SSH Server" > /dev/stdout

exec sudo /usr/sbin/sshd -f /ssh/sshd_config/config -D -e
