# etc/s6-overlay/s6-rc.d/fix-timezone/up
#!/command/execlineb -P

/command/with-contenv sudo /var/fixes/set_timezone.sh
