# etc/s6-overlay/s6-rc.d/fix-project-symlinks/up
#!/command/execlineb -P

/command/with-contenv /var/fixes/set_symlinks.sh
